# deployment/systemd/amazon-notify.service
[Unit]
Description=Amazon Notify (Gmail -> Discord)
After=network-online.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/opt/amazon-notify-v2
ExecStart=/opt/amazon-notify-v2/.venv/bin/amazon-notify --config /opt/amazon-notify-v2/config.json
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
