# deployment/systemd/amazon-notify-fallback.timer
[Unit]
Description=Amazon Notify fallback poll timer

[Timer]
OnBootSec=2m
OnUnitActiveSec=10m
RandomizedDelaySec=60
Unit=amazon-notify-fallback.service
Persistent=true

[Install]
WantedBy=timers.target
