# deployment/systemd/amazon-notify.service
[Unit]
Description=Amazon Notify (Gmail -> Discord)
After=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=3
OnFailure=amazon-notify-alert@%n.service

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/opt/amazon-notify
ExecStart=/opt/amazon-notify/.venv/bin/amazon-notify --config /opt/amazon-notify/config.json
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
