# deployment/systemd/amazon-notify-fallback.service
[Unit]
Description=Amazon Notify fallback poll (oneshot)
After=network-online.target

[Service]
Type=oneshot
User=YOUR_USER
WorkingDirectory=/opt/amazon-notify
ExecStart=/opt/amazon-notify/.venv/bin/amazon-notify --config /opt/amazon-notify/config.json --once --fallback-watchdog --main-service-name amazon-notify-pubsub.service --heartbeat-file /opt/amazon-notify/runtime/pubsub-heartbeat.txt --heartbeat-max-age-seconds 120
